# compose/app/app.ini
[unix_http_server]
file=/tmp/supervisor.sock    ; the path to the socket file

[supervisord]
logfile=/var/log/app/supervisord.log
logfile_maxbytes=1000MB
logfile_backups=10
loglevel=debug
pidfile=/tmp/supervisord.pid
nodaemon=true
childlogdir=/var/log/app

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:app-masterd]
user=root
directory=/app
command=python audicus/app.py
